# agent/scripts/com.bumba.agent-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bumba.agent-maintenance</string>
    <key>UserName</key>
    <string>bumba-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/opt/bumba-harness/agent-flat/agent/scripts/maintenance.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/opt/bumba-harness/agent-flat/agent</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/opt/bumba-harness/logs/maintenance-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/bumba-harness/logs/maintenance-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/opt/bumba-harness</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
