# agent/config/launchdaemons/com.bumba.agent-experiment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bumba.agent-experiment</string>

    <key>UserName</key>
    <string>bumba-agent</string>

    <key>WorkingDirectory</key>
    <string>/opt/bumba-harness/agent-flat/agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/bumba-harness/agent-flat/agent/.venv/bin/python3</string>
        <string>scripts/experiment_loop.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/opt/bumba-harness</string>
        <key>PATH</key>
        <string>/opt/bumba-harness/.local/bin:/opt/bumba-harness/agent-flat/agent/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/opt/bumba-harness/logs/experiment-loop-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/bumba-harness/logs/experiment-loop-stderr.log</string>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
